# .config/systemd/user/keynav.service
[Unit]
Description=launch keynav

[Service]
ExecStart=/usr/bin/keynav
Restart=always

[Install]
WantedBy=xorg.target
