# rootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStartPre=/bin/sudo -ubrnv /bin/touch /tmp/autologin
ExecStart=-/usr/bin/agetty --autologin brnv --noclear %I 38400 linux
